# rootfs/etc/s6-overlay/s6-rc.d/nginx/finish
#!/command/execlineb -S0
# Nginx finish script - auto-restart on crash
if { s6-test ${1} -ne 0 }
if { s6-test ${1} -ne 256 }
s6-svc -r .
