# rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/bin/sh
# S6 dependency ensures init-configs runs before this
echo "[Nginx] Starting nginx (configs ready via S6 dependency)..."
exec nginx -g "daemon off;"
