# rootfs/etc/s6-overlay/s6-rc.d/nginx/up
#!/command/execlineb -P
# Wait for nginx.conf to exist and be valid before starting
foreground { sleep 2 }
if { test -f /etc/nginx/nginx.conf }
foreground { nginx -t }
exit 0
