# rootfs/etc/s6-overlay/s6-rc.d/nginx-config-ready/up
#!/bin/bash
# Wait for nginx configs to be generated (max 10 seconds)
for i in {1..50}; do
    if [ -f /etc/nginx/nginx.conf ] && [ -f /etc/nginx/http.d/default.conf ]; then
        echo "[S6] Nginx configs ready, validating..."
        nginx -t
        exit 0
    fi
    sleep 0.2
done

echo "[S6] ERROR: Nginx configs not found after 10 seconds!"
exit 1
